# lumtek-web/deploy/lumtek-api.service.example
# Copiar a /etc/systemd/system/lumtek-api.service
# systemctl daemon-reload && systemctl enable --now lumtek-api

[Unit]
Description=Lumtek contact API (Node + SMTP)
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/lumtek-app
EnvironmentFile=/var/www/lumtek-app/.env
ExecStart=/usr/bin/node server/index.mjs
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
